# deploy/fusion-ocr-worker.service
# systemd unit for the fusion-ocr queue WORKER (Linux hosts).
# The API (fusion-ocr.service) only enqueues; this drains the queue and runs the pipeline.
# Run BOTH. They share in/, out/ and out/jobs.sqlite.
#   sudo cp deploy/fusion-ocr-worker.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now fusion-ocr-worker

[Unit]
Description=fusion-ocr queue worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=fusion-ocr
WorkingDirectory=/opt/fusion-ocr
Environment=PYTHONUNBUFFERED=1   # progress/status lines reach the journal live, not block-buffered
# Same env as the API (airgap pre-pulled models, HF_TOKEN, etc.). The worker does NOT need
# FUSION_OCR_API_TOKEN — it's not an HTTP service.
EnvironmentFile=/etc/fusion-ocr/fusion-ocr.env

ExecStart=/opt/fusion-ocr/.venv/bin/fusion-ocr        # the watcher loop = the worker
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/fusion-ocr/in /opt/fusion-ocr/out

[Install]
WantedBy=multi-user.target
